# infra/supervisord.conf
# =============================================================================
# supervisord.conf — Amprealize API (Fly.io production)
# =============================================================================
# Manages two processes inside the single Fly.io container:
#   1. nginx  — gateway on :8080  (rate limiting, header stripping, real-IP)
#   2. uvicorn — FastAPI on 127.0.0.1:8000  (internal only, never exposed)
#
# Start order: nginx waits for uvicorn to be healthy before accepting traffic
# because the Fly.io health check hits /health through nginx.  Both processes
# auto-restart on failure; supervisord exits if either crashes repeatedly.
# =============================================================================

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid

# -------------------------------------------------------------------------
# nginx gateway — starts first, proxies to uvicorn on localhost
# -------------------------------------------------------------------------
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# -------------------------------------------------------------------------
# uvicorn (FastAPI) — binds only on loopback, not reachable externally
# -------------------------------------------------------------------------
# --proxy-headers: honour X-Real-IP / X-Forwarded-For set by nginx.
# --forwarded-allow-ips 127.0.0.1: only trust proxy headers from nginx
#   (not from the public internet or Fly.io private network).
[program:uvicorn]
command=uvicorn amprealize.api:app
    --host 127.0.0.1
    --port 8000
    --workers 1
    --proxy-headers
    --forwarded-allow-ips 127.0.0.1
user=amprealize
directory=/app
autostart=true
autorestart=true
priority=20
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/app",USER="amprealize",PYTHONUNBUFFERED="1"
